# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:haproxy]
command=/usr/sbin/haproxy -f /etc/haproxy/haproxy.cfg

[program:php5-fpm]
command = /usr/sbin/php7-fpm

[program:nginx]
command=/usr/sbin/nginx
